# raspberry pi/relay/ac_script.service
[Unit]
Description=run AC

[Service]
ExecStart=/usr/bin/bash /usr/bin/ac_script.sh

[Install]
WantedBy=multi-user.target